# roles/hypr/files/change-wallpaper.timer
[Timer]
OnBootSec=60
OnCalendar=hourly

[Install]
WantedBy=graphical-session.target
